# prometheus_dcgm_exporter/templates/dcgm_exporter.service.j2
[Unit]
Description=DCGM Exporter
After=network.target
[Service]
User=nodeusr
Group=nodeusr
Type=simple
ExecStart={{ executable_path }}/dcgm_exporter
[Install]
WantedBy=multi-user.target
